# userpatches/overlay/etc/systemd/system/setup-containers.service
[Unit]
Description=Setup Docker Containers on First Boot
After=docker.service
Requires=docker.service
ConditionPathExists=!/var/lib/setup-containers-done

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-portainer.sh
ExecStart=/usr/local/bin/setup-octoprint.sh
ExecStartPost=/usr/bin/touch /var/lib/setup-containers-done
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
